# config/fingerprint/enroll-fingerprint.service
[Unit]
Description=One-time fingerprint enrollment prompt
After=graphical.target

[Service]
Type=oneshot
ExecStart=/usr/bin/kitty --hold -e /home/james/Dev/hyprland/scripts/enroll_fingerprint.sh
Environment=DISPLAY=:1
Environment=XDG_SESSION_TYPE=x11
Environment=XAUTHORITY=%h/.Xauthority

[Install]
WantedBy=default.target
