# config/fingerprint/enroll-fingerprint.service
[Unit]
Description=One-time fingerprint enrollment prompt
After=graphical.target

[Service]
Type=oneshot
ExecStart=/usr/bin/kitty --hold -e /home/james/Dev/hyprland/scripts/enroll_fingerprint.sh

[Install]
WantedBy=default.target
